# device-startup-0.0.2/debian/postinst
#!/bin/sh
update-rc.d wifi defaults
update-rc.d usb defaults
#DEBHELPER#
