# bananui-base-0.1.1/debian/postrm
#!/bin/sh
set -e
if [ "$1" = "purge" ] ; then
	update-rc.d bananui remove >/dev/null
fi
#DEBHELPER#
