# bananui-base-0.1.1/debian/postinst
#!/bin/sh
update-rc.d bananui defaults
#DEBHELPER#
