# bananui-base-0.1.1/debian/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/bananui" ]; then
	invoke-rc.d bananui stop
fi
#DEBHELPER#
