# bananui-base-0.0.2/debian/postrm
#!/bin/sh
update-rc.d bananui remove
#DEBHELPER#
