# device-startup-0.1.1/debian/postrm
#!/bin/sh
update-rc.d wifi remove
update-rc.d usb remove
#DEBHELPER#
